# gnu_dataset/gnats-4.2.0/debian/gnats.cron.daily
#!/bin/sh

su - gnats -c 'test -x /usr/lib/gnats/check-db && /usr/lib/gnats/check-db --all-databases; exit 0'
