# gnu_dataset/sauce-0.7.7/debian/cron.weekly
#!/bin/sh
# $Id: cron.weekly,v 1.2 1998/09/02 19:08:14 ian Exp $

. /etc/sauce/sys-config

test -x /usr/sbin/sauceadmin || exit 0
egrep '^\#.*sauce' /etc/inittab >/dev/null || exit 0

# Cycle logs

cd /var/log/sauce
savelog -c $SAUCE_LOG_SAVE -m $SAUCE_LOG_MODE -u $SAUCE_USER -g $SAUCE_GROUP \
	startup.log debug.log notice.log dbreasons.log error.log fatal.log
sauceadmin reopenlogs
